# mrmShared/linux/dkms/50-mrf.rules.in
# UDEV rule to give general access to MRF device
SUBSYSTEM=="uio", ATTR{name}=="mrf-pci", GROUP="@@MRF_USER@@", MODE="0660", RUN+="/bin/chmod g+rw %S%p/resource0"
